# supervisor/rsyslogd.ini
[program:rsyslogd]
command=/usr/sbin/rsyslogd -n -c5
redirect_stderr=true
